# scripts/crontab/auto_change_wallpaper.sh
#!/bin/bash

# Load environments
export DISPLAY=:0

/usr/bin/awsetbg -r /home/mudox/.bg_papers >> /home/mudox/.log/cron/auto_wall_paper.log 2>&1

